# robot-pi/miit-rover-agent.service
[Unit]
Description=MIIT Rover MQTT bridge
Wants=network-online.target time-sync.target
After=network-online.target time-sync.target
ConditionPathExists=/opt/miit-rover/source/robot-pi/agent.py
ConditionPathExists=/opt/miit-rover/source/robot-pi/message_contract.py
ConditionPathExists=/opt/miit-rover/source/robot-pi/local_store.py

[Service]
Type=simple
User=rover
Group=rover
SupplementaryGroups=dialout
WorkingDirectory=/opt/miit-rover/source/robot-pi
EnvironmentFile=/etc/miit-rover/robot.env
ExecStart=/opt/miit-rover/source/robot-pi/.venv/bin/python /opt/miit-rover/source/robot-pi/agent.py
Restart=always
RestartSec=3
TimeoutStopSec=10
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/miit-rover

[Install]
WantedBy=multi-user.target
